# infra/mcp-servers/setup/dedup-graph.service
# Server-side dedup: collapse case-variant / drifted duplicate slugs that branch
# merges can introduce (Omnigraph only auto-merges SAME-slug nodes). Idempotent
# and a cheap no-op when the graph is clean; rebuilds only when duplicates exist.
# Runs on the memory SERVER host (coding.vm), not on clients.
# Install: cp to ~/.config/systemd/user/, then: systemctl --user enable --now dedup-graph.timer
[Unit]
Description=Deduplicate the Omnigraph memory graph
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/code/agent-skills/infra/mcp-servers/scripts/dedup-graph.py
Nice=15
