# infra/mcp-servers/setup/omnigraph-sync.service
# Install: copy this + omnigraph-sync.timer to ~/.config/systemd/user/ (or
# /etc/systemd/system/ for a system service), set ExecStart to the script path,
# then: systemctl --user enable --now omnigraph-sync.timer
[Unit]
Description=Reconcile local Omnigraph memory with the central server
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Point at your checkout; secrets/config resolved from setup/.env next to it.
ExecStart=%h/code/agent-skills/infra/mcp-servers/setup/omnigraph-sync.sh
Nice=10
